# systemd/tria-aibpo-backup.timer
# Tria AIBPO backup timer
# Schedules automated backups daily at 2 AM
#
# Installation:
#   sudo cp systemd/tria-aibpo-backup.* /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tria-aibpo-backup.timer
#   sudo systemctl start tria-aibpo-backup.timer
#
# Check timer status:
#   sudo systemctl status tria-aibpo-backup.timer
#   sudo systemctl list-timers tria-aibpo-backup.timer

[Unit]
Description=Tria AIBPO Daily Backup Timer
Documentation=https://github.com/yourusername/tria-aibpo
Requires=tria-aibpo-backup.service

[Timer]
# Run daily at 2 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run 15 minutes after boot (if system was off during scheduled time)
OnBootSec=15min

# If a backup was missed (system was off), run it when system starts
Persistent=true

# Accuracy (allow up to 15 min deviation for power saving)
AccuracySec=15min

[Install]
WantedBy=timers.target
